# baby-sweep/baby-sweep.service
[Unit]
Description=Baby Sweep Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/baby-sweep/backend
Environment="NODE_ENV=production"
Environment="PORT=3001"
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
